# check-hdfs-space.service
[Unit]
Description=send hdfs remaining space to zabbix

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/check-hdfs-space.sh